# services/claude-proxy/jobclub-claude-proxy.service
[Unit]
Description=Job Club Claude proxy (shells out to `claude` CLI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jobclub-claude-proxy
EnvironmentFile=/opt/jobclub-claude-proxy/.env
ExecStart=/opt/jobclub-claude-proxy/venv/bin/gunicorn \
    --bind 0.0.0.0:8090 \
    --workers 2 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
